# flask-app.service
[Unit]
Description=Flask app service
After=network.target

[Service]
WorkingDirectory=/home/mixola/app
ExecStart=/home/mixola/app/venv/bin/python main.py

[install]
WantedBy=multi-user.target